# .config/systemd/user/batterymon.service
[Unit]
Description=Batterymon
After=xorg.target

[Service]
ExecStart=/usr/bin/batterymon

[Install]
WantedBy=xinitrc.target
