# .config/systemd/user/dunst.service
[Unit]
Description=Dunst notification daemon
After=xorg.target

[Service]
ExecStart=/usr/bin/dunst
Type=forking

[Install]
WantedBy=xinitrc.target
