# .config/systemd/user/gpodder.timer
[Unit]
Description=Fetch new podcast episodes each day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
