# .config/systemd/user/maildirwatch.service
[Unit]
Description=Maildir notifier
After=xorg.target

[Service]
ExecStart=/home/erikw/bin/maildir_watch.sh

[Install]
WantedBy=xinitrc.target
