# .config/systemd/user/netmon.service
[Unit]
Description=Netmon network monitor for netctl
After=xorg.target

[Service]
ExecStart=/usr/bin/netmon

[Install]
WantedBy=xinitrc.target
