# .config/systemd/user/notmuch.service
[Unit]
Description=Index email with notmuch

[Service]
Type=simple
Nice=10
ExecStart=/usr/bin/notmuch new
