# .config/systemd/user/parcellite.service
[Unit]
Description=Parcellite
After=xorg.target

[Service]
ExecStart=/usr/bin/parcellite

[Install]
WantedBy=xinitrc.target
